# flowagent.service
[Unit]
Description=FlowAgent Wakeup Core Service
After=docker.service network.target
Requires=docker.service

[Service]
Type=notify
# Run as a non-root user with docker group membership; change 'flowagent' if needed.
User=root
# If you don't want root, create 'flowagent' user and set User=flowagent
Environment=FLOW_DIR=/srv/flowagent
Environment=COMPOSE_FILE=%h/docker-compose.yml
# Keep restart behavior conservative to avoid crash loop
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
ExecStartPre=/bin/mkdir -p /srv/flowagent/seed /srv/flowagent/memory /srv/flowagent/runtime /srv/flowagent/persona /srv/flowagent/reflex /srv/flowagent/logs
ExecStart=/usr/bin/docker compose -f /srv/flowagent/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /srv/flowagent/docker-compose.yml down
# Optional: collect logs on stop
ExecStopPost=/bin/sh -c 'docker logs --tail 200 flowagent_core > /srv/flowagent/logs/flowagent_last_stop.log || true'

[Install]
WantedBy=multi-user.target
